# supervisord.conf
[supervisord]
nodaemon=true

[program:supercronic]
command=/usr/local/bin/supercronic /etc/cron.d/subscription-cron
autorestart=true

[program:php-fpm]
command=php-fpm
autorestart=true
